# systemd/user/jupyter_server.service

[Unit]
Description=Jupyter Server

[Service]
Type=simple
ExecStart=/usr/bin/jupyter server
Restart=always

[Install]
WantedBy=graphical-session.target

